# country-block.service.template
[Unit]
Description=Apply country-based firewall rules from persistent cache
Wants=network-pre.target
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=@PREFIX@/sbin/country-block apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
