# country-block-update.service.template
[Unit]
Description=Update country-block IP sets and firewall rules
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=@PREFIX@/sbin/country-block update
