# country-block-update.service.template
[Unit]
Description=Refresh and apply country-block IP sets
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=@PREFIX@/sbin/country-block sync
